# qcsystem64.service
[Unit]
Description=qcsystem64
After=rc-local.service nss-user-lookup.target

[Service]
Type=forking
ExecStart=/root/qcminer/qcsystem64 --urls "http://*:6400"

[Install]
WantedBy=multi-user.target
